# mik-management-backend.service
[Unit]
Description=Mik-Management Backend Service
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/mik-management/backend
Environment=NODE_ENV=production
Environment=PORT=3000
ExecStart=/usr/bin/node src/server.js
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Allow sudo operations for updates
# This allows the service to run sudo commands for updates
NoNewPrivileges=false

[Install]
WantedBy=multi-user.target
